# service/keystone/run
#!/bin/sh

exec 2>&1

exec runuser -u keystone -g keystone -- /usr/bin/keystone-all

